# root/services/cron.conf
[program:cron]
command=/usr/bin/crond -n
