# root/etc/crontab/root
0 3 * * * /etc/cont-finish.d/backup-db.sh
